# idmind_robot/no_build/config/udev_rules/11-projector.rules
KERNEL=="card0", SUBSYSTEM=="drm", ACTION=="change", RUN+="/home/viva/catkin_ws/src/idmind_robot/no_build/scripts/projector.sh"
